# BHCOL/SCRIPTS/BHCOL.cron
05 * * * * /opt/BHCOL/SCRIPTS/bh_collect.sh
